# 3_network_connection_initiated/include_living_of_the_land.xml
<Sysmon schemaversion="4.3">
   <EventFiltering>
 <RuleGroup name="" groupRelation="or">
      <NetworkConnect onmatch="include">
	  		<!--Native Windows tools - Living off the land-->
			<Image name="technique_id=T1016,technique_name=System Network Configuration Discovery" condition="image">ipconfig.exe</Image><!--Microsoft:Windows: shows ip configuration -->
			<Image name="technique_id=T1057,technique_name=Process Discovery" condition="image">tasklist.exe</Image> <!--Mitre T1057--><!--Mitre T1007--><!--Microsoft:Windows: shows current running processes-->
			<Image name="technique_id=T1049,technique_name=System Network Connections Discovery" condition="image">netstat.exe</Image> <!--Mitre T1016--><!--Mitre T1049--><!--Microsoft:Windows: shows protocol statistics and current TCP/IP network connections -->
			<Image name="technique_id=T1057,technique_name=Process Discovery" condition="image">qprocess.exe</Image> <!--Microsoft:Windows: shows information about processes -->
      <Image name="technique_id=T1016,technique_name=System Network Configuration Discovery" condition="image">nslookup.exe</Image><!--Microsoft:Windows: shows DNS configuration and enables quering -->
      <Image name="technique_id=T1018,technique_name=Remote System Discovery" condition="image">net.exe</Image> <!--Mitre T1201--><!--Mitre T1007--><!--Mitre T1124--><!--Mitre T1018--><!--Mitre T1049--><!--Mitre T1087--><!--Mitre T1135--><!--Mitre T1069--><!--Mitre T1016--><!--Mitre T1077--><!--Microsoft:Windows: "net use"/"net view" used by attackers to surveil and connect with file shares from command line | Credit @ion-storm -->
			<Image name="technique_id=T1033,technique_name=System Owner/User Discovery" condition="image">quser.exe</Image> <!--Microsoft:Windows: shows logged-on users -->
			<Image name="technique_id=T1057,technique_name=Process Discovery" condition="image">query.exe</Image> <!--Microsoft:Windows: shows information about processes -->
			<Image name="technique_id=T1134,technique_name=Access Token Manipulation" condition="image">runas.exe</Image> <!--Microsoft:Windows: run a process as another user -->
			<Image name="technique_id=T1112,technique_name=Modify Registry" condition="image">reg.exe</Image> <!--Mitre T1214--><!--Mitre T1012--><!--Mitre T1112--><!--Microsoft:Windows: reads and modifies the Windows register -->
			<Image name="technique_id=T1063,technique_name=Security Software Discovery" condition="image">netsh.exe</Image> <!--Mitre T1016--><!--Mitre T1128--><!--Mitre T1063--><!--Microsoft:Windows: manipulate the firewall -->
			<Image condition="image">klist.exe</Image> <!--Microsoft:Windows: show cached kerberos tickets -->
			<Image name="technique_id=T1070,technique_name=Indicator Removal on Host" condition="image">wevtutil.exe</Image><!--Microsoft:Windows: read and modify the Windows Eventlog -->
      <Image name="technique_id=T1053,technique_name=Scheduled Task" condition="image">taskeng.exe</Image><!--Microsoft:Windows: taskscheduler -->
			<Image name="technique_id=T1117,technique_name=Regsvr32" condition="image">regsvr32.exe</Image> <!--Mitre T1117--><!--Mitre T1085--><!--Microsoft:Windows: [ https://subt0x10.blogspot.com/2016/04/bypass-application-whitelisting-script.html ] -->
      <Image condition="image">hh.exe</Image> <!--Microsoft:Windows: HTML Helper-->
      <Image name="technique_id=T1059,technique_name=Command-Line Interface" condition="image">cmd.exe</Image><!--Microsoft:Windows: Command prompt-->
      <Image name="technique_id=T1086,technique_name=PowerShell" condition="image">powershell.exe</Image><!--Microsoft:Windows: PowerShell interface-->
      <Image name="technique_id=T1202,technique_name=Indirect Command Execution" condition="image">bash.exe</Image><!--Microsoft:Windows: bash on Windows, Linux subsystem-->
      <Image name="technique_id=T1202,technique_name=Indirect Command Execution" condition="image">pcalua.exe</Image><!--Microsoft:Windows: Program Compatibility Assistant)-->
      <Image name="technique_id=T1202,technique_name=Indirect Command Execution" condition="image">cscript.exe</Image>
      <Image name="technique_id=T1202,technique_name=Indirect Command Execution" condition="image">wscript.exe</Image>
      <Image name="technique_id=T1170,technique_name=Mshta" condition="image">mshta.exe</Image>
      <Image name="technique_id=T1016,technique_name=System Network Configuration Discovery" condition="image">nbtstat.exe</Image> <!--Microsoft:Windows: NetBIOS statistics, attackers use to enumerate local network -->
      <Image name="technique_id=T1069,technique_name=Permission Groups Discovery" condition="image">net1.exe</Image> <!--Mitre T1018--><!--Mitre T1077--><!--Mitre T1087--><!--Mitre T1135--><!--Mitre T1069--><!--Mitre T1016--><!--Microsoft:Windows: Launched by "net.exe", but it may not detect connections either -->
      <Image name="technique_id=T1018,technique_name=Remote System Discovery" condition="image">nslookup.exe</Image> <!--Microsoft:Windows: Retrieve data over DNS -->
      <Image name="technique_id=T1057,technique_name=Process Discovery" condition="image">qwinsta.exe</Image> <!--Microsoft:Windows: Remotely query login sessions on a server or workstation | Credit @ion-storm -->
      <Image name="technique_id=T1057,technique_name=Process Discovery" condition="image">rwinsta.exe</Image> <!--Microsoft:Windows: Disconnect remote sessions | Credit @ion-storm -->
      <Image condition="image">schtasks.exe</Image> <!--Microsoft:Windows: Command-line interface to local and remote tasks -->
      <Image condition="image">taskkill.exe</Image> <!--Microsoft:Windows: Kill processes, has remote ability -->
      <Image name="technique_id=T1031,technique_name=Modify Existing Service" condition="image">sc.exe</Image>
      <Image name="technique_id=T1033,technique_name=System Owner/User Discovery" condition="image">nltest.exe</Image><!-- netlogon testing utility-->
      <Image name="technique_id=T1202,technique_name=Indirect Command Execution" condition="image">winrs.exe</Image>
		</NetworkConnect>
</RuleGroup>
</EventFiltering>
</Sysmon>
